# ztest_lab1/monitor_etc.service.j2
[Unit]
Description=Monitor /etc directory

[Service]
Type=simple
ExecStart=/usr/local/bin/monitor_etc.sh

[Install]
WantedBy=multi-user.target
